# rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty 115200 console
::respawn:syslogd -n
